# systemd/anthropic-lb.service
[Unit]
Description=Anthropic API Load Balancer
After=network.target

[Service]
Type=simple
# Adjust User/Group and paths for your deployment
User=user
Group=user
WorkingDirectory=/home/user/anthropic-load-balancer
Environment="PATH=/home/user/anthropic-load-balancer/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/user/anthropic-load-balancer/venv/bin/uvicorn src.main:app --host 0.0.0.0 --port 8080 --workers 1
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/user/anthropic-load-balancer

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anthropic-lb

[Install]
WantedBy=multi-user.target
